# deploy/vrab-engine.service
[Unit]
Description=VRAB Trading Engine
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/vrab
EnvironmentFile=/root/vrab/.env
ExecStart=/root/vrab/venv/bin/python -m live.engine
Restart=on-failure
RestartSec=10
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
