# deploy/vrab-dashboard.service
[Unit]
Description=VRAB Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/vrab
EnvironmentFile=/root/vrab/.env
ExecStart=/root/vrab/venv/bin/python -m dashboard
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
